# Raspberry_Pi/packages/mediamtx_1.15.0_arm64/DEBIAN/prerm
#!/bin/sh -e

set -e

if [ "$1" = "remove" ]; then
	systemctl stop mediamtx.service
	systemctl disable mediamtx.service
fi
