# Raspberry_Pi/packages/mediamtx_1.15.0_arm64/lib/systemd/system/mediamtx.service
[Unit]
Wants=network.target
[Service]
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/mediamtx /etc/mediamtx.yml
[Install]
WantedBy=multi-user.target
